# automation/services/mina-bp-stats/leaderboard-bot/minanet_app/crontab
*/2 * * * * /usr/local/bin/python /opt/minanet/survey_collect.py 

